# .config/autostart/_autostart_foot_server.desktop
[Desktop Entry]
Exec=foot --server
Icon=dialog-scripts
Name=server for foo terminal
Path=
Type=Application
X-KDE-AutostartScript=true
