# .config/autostart/_st_login.sh.desktop
[Desktop Entry]
Exec=/home/st/.local/bin/_st_login.sh
Icon=dialog-scripts
Name=_st_login.sh
Path=
Type=Application
X-KDE-AutostartScript=true
