# .config/autostart/unclutter.desktop
[Desktop Entry]
Exec=unclutter --timeout 1
Icon=
Name=unclutter
Path=
Terminal=False
Type=Application
