# .config/autostart/wwconsole.desktop
[Desktop Entry]
Exec=/home/st/.local/bin/wwconsole
Icon=dialog-scripts
Name=wwconsole
Path=
Type=Application
X-KDE-AutostartScript=true
